# examples/systemd/vdo.service
[Unit]
Description=VDO volume services (provided for compatibility)
Documentation=man:vdo(8)
